# crontab
19 5 * * * root /usr/local/bin/nicedocker exec $service_id /usr/sbin/logrotate /etc/logrotate.d/redis-server
