# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
